# deploy/trading-system.service
[Unit]
Description=APEX Multi-Agent Trading System
After=ibgateway.service network-online.target
Wants=ibgateway.service network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/agentic-trader
Environment=PATH=/root/agentic-trader/venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=DISPLAY=:1
EnvironmentFile=/root/agentic-trader/.env

# Full trading system (manual approval mode - no --auto flag)
ExecStart=/root/agentic-trader/venv/bin/python main.py

# Auto-restart on failure
Restart=always
RestartSec=15

# Wait for IB Gateway to be ready before starting
ExecStartPre=/bin/sleep 30

# Logging
StandardOutput=journal
StandardError=journal

# Resource limits
MemoryLimit=1G
CPUQuota=80%

[Install]
WantedBy=multi-user.target
